# weekly/root_overlays/openstack/etc/profile.d/tcc.sh
export CC="tcc"
export LDSHARED="tcc -shared"
