# weekly/root_overlays/steamcmd/etc/systemd/system/steamcmd@294420.service.d/override.conf
[Unit]
Description=7 Days to Die Dedicated Server

[Install]
Alias=7days@294420.service

[Service]
ExecStart=/opt/steamcmd/294420/startserver.sh -configfile=serverconfig.xml
