# weekly/root_overlays/steamcmd/etc/systemd/system/steamcmd@1690800.service.d/override.conf
[Unit]
Description=Satisfactory Dedicated Server

[Install]
Alias=satisfactory@1690800.service

[Service]
ExecStart=/opt/steamcmd/1690800/FactoryServer.sh
