# weekly/root_overlays/openstack/etc/systemd/system/rabbitmq.service.d/override.conf
[Service]
StateDirectory=%p
LogsDirectory=%p
